# backend/scripts/crontab
# Run backup daily at 2 AM

0 2 * * * /backup.sh >> /backups/backup.log 2>&1
