# backend/scripts/crontab
0 2 * * * root /backup.sh >> /var/log/backup.log 2>&1