# scripts/crontab
# Run backup daily at 2 AM

0 2 * * * /srv/cortexex/scripts/backup_docker.sh >> /srv/cortexex/backups/backup.log 2>&1